# templates/my-docker-image.service.j2
[Unit]
Description=My docker image; Final work; Task 1
Requires=docker.service

[Service]
User=root
ExecStart=/bin/sh -c "sudo docker run {{ docker_image }}"

[Install]
WantedBy=multi-user.target

